# packaging/dev.staywake.daemon.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.staywake.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>__STAYWAKE_BIN__</string>
    <string>daemon</string>
    <string>--state-path</string>
    <string>__STATE_PATH__</string>
    <string>--config-path</string>
    <string>__CONFIG_PATH__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/var/log/staywake.log</string>

  <key>StandardErrorPath</key>
  <string>/var/log/staywake.log</string>
</dict>
</plist>
